# reconcile
#! /bin/sh
### BEGIN INIT INFO
# Provides:          reconcile
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Controls reconcile
# Description:       An init script for reconcile.
### END INIT INFO

# Author: Roelf Wichertjes <roelf@roelf.org>

# CONFIG

BOTUSER=reconcile
BOTDIR=/home/$BOTUSER/Reconcile
# Add -3 to this if you have a python3 instead of python bin.
# It still needs the python version to be 3 or higher though.
# Please don't remove -d. Your system might freeze if you do.
BOTARGS="-d"

# END CONFIG

# DO NOT EDIT BEYOND THIS POINT

# Do NOT "set -e"

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin
NAME="reconcile"
DAEMON=$(ls $BOTDIR/start.sh)
SCRIPTNAME=/etc/init.d/$NAME

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
	WD=$(pwd)
	cd $BOTDIR
	sudo -u $BOTUSER ./start.sh $BOTARGS
	RETVAL="$?"
	cd $WD
	return "$RETVAL"
}

#
# Function that stops the daemon/service
#
do_stop()
{
	WD=$(pwd)
	cd $BOTDIR
	sudo -u $BOTUSER ./start.sh -k
	RETVAL="$?"
	cd $WD
	return "$RETVAL"
}

case "$1" in
  start)
	[ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESC" "$NAME"
	do_start
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	;;
  stop)
	[ "$VERBOSE" != no ] && log_daemon_msg "Stopping $DESC" "$NAME"
	do_stop
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	;;
  restart|force-reload)
	[ "$VERBOSE" != no ] && log_daemon_msg "Stopping $DESC" "$NAME"
	do_stop
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	[ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESC" "$NAME"
	do_start
	case "$?" in
		0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
		2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
	esac
	;;
  *)
	echo "Usage: $SCRIPTNAME {start|stop|restart|force-reload}" >&2
	exit 3
	;;
esac


